# crontab
0 */2 * * * bash /opt/rclone/scripts/upcall
